# scripts/nixx-server.service
[Unit]
Description=nixx API server
After=network.target docker.service nixx-embed.service
Requires=docker.service
Wants=nixx-embed.service

[Service]
Type=simple
User=siderealyear
WorkingDirectory=/home/siderealyear/nixx
ExecStart=/home/siderealyear/nixx/venv/bin/nixx serve
Restart=on-failure
RestartSec=5s
EnvironmentFile=/home/siderealyear/nixx/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=nixx.target
